# tests/golden/render/freebsd/example-local-dns-redirect.golden
### dnsmasq.conf
# Generated by routerd. Do not edit by hand.
dhcp-leasefile=/var/db/routerd/dnsmasq/dnsmasq.leases
port=0
no-resolv
bind-dynamic
interface=ens19

interface=ens19
dhcp-range=set:lan-dhcpv4,192.168.50.100,192.168.50.199,12h
dhcp-option=tag:lan-dhcpv4,option:dns-server,192.168.50.1
dhcp-option=tag:lan-dhcpv4,option:domain-name,home.example
dhcp-option=tag:lan-dhcpv4,option:domain-search,home.example
### pf.conf
# Generated by routerd. Do not edit by hand.
set skip on lo0
lan_if = "ens19"
block drop all
pass quick on lo0 all
pass out quick all keep state
pass quick inet6 proto icmp6 all keep state
pass in quick on $lan_if to self keep state
pass in quick on $lan_if keep state label "routerd:lan-to-lan"
### rc.conf.d-routerd
# Generated by routerd. Do not edit by hand.
gateway_enable="YES"
ipv6_gateway_enable="YES"
ifconfig_ens19="inet 192.168.50.1/24"
pf_enable="YES"
pflog_enable="YES"
### rc.d-routerd_dnsmasq
#!/bin/sh

# PROVIDE: routerd_dnsmasq
# REQUIRE: NETWORKING mpd5
# KEYWORD: shutdown

. /etc/rc.subr

name="routerd_dnsmasq"
rcvar="routerd_dnsmasq_enable"
command="/usr/local/sbin/dnsmasq"
pidfile="/var/run/routerd/dnsmasq.pid"
command_args="--conf-file=/usr/local/etc/routerd/dnsmasq.conf --pid-file=${pidfile}"
start_precmd="${name}_prestart"

routerd_dnsmasq_prestart()
{
	mkdir -p "/var/run/routerd"
	mkdir -p "/var/db/routerd/dnsmasq"
}

load_rc_config ${name}
: ${routerd_dnsmasq_enable:="NO"}

run_rc_command "$1"
