# tests/golden/render/freebsd/freebsd-vrrp.golden
### pf.conf
# Generated by routerd. Do not edit by hand.
set skip on lo0
lan_if = "vtnet1"
block drop all
pass quick on lo0 all
pass out quick all keep state
pass quick inet6 proto icmp6 all keep state
pass in quick on $lan_if to self keep state
pass in quick on $lan_if keep state label "routerd:lan-to-lan"
### rc.conf.d-routerd
# Generated by routerd. Do not edit by hand.
gateway_enable="YES"
ipv6_gateway_enable="YES"
ifconfig_vtnet1="inet 192.168.70.2/24"
pf_enable="YES"
pflog_enable="YES"
routerd_carp_enable="YES"
### rc.d-routerd_carp
#!/bin/sh
#
# Managed by routerd.
#
# PROVIDE: routerd_carp
# REQUIRE: NETWORKING
# KEYWORD: shutdown

. /etc/rc.subr

name="routerd_carp"
rcvar="routerd_carp_enable"
start_cmd="routerd_carp_start"
stop_cmd="routerd_carp_stop"
status_cmd="routerd_carp_status"

routerd_carp_start() {
  kldload carp >/dev/null 2>&1 || true
  sysctl net.inet.carp.preempt='0' >/dev/null
  ifconfig 'vtnet1' 'inet' 'vhid' '70' 'advbase' '1' 'advskew' '104' 'pass' 'change-me' 'alias' '192.168.70.10/32'
}

routerd_carp_stop() {
  ifconfig 'vtnet1' inet '192.168.70.10/32' -alias >/dev/null 2>&1 || true
}

routerd_carp_status() {
  true && ifconfig 'vtnet1' | grep -q 'vhid 70'
}

load_rc_config $name
: ${routerd_carp_enable:="YES"}
run_rc_command "$1"
