# tests/golden/render/freebsd/tailscale-minimal.golden
### rc.conf.d-routerd
# Generated by routerd. Do not edit by hand.
gateway_enable="YES"
ipv6_gateway_enable="YES"
tailscaled_enable="YES"
routerd_tailscale_mesh_enable="YES"
### rc.d-routerd_tailscale_mesh
#!/bin/sh
#
# Managed by routerd.
#
# PROVIDE: routerd_tailscale_mesh
# REQUIRE: NETWORKING tailscaled
# KEYWORD: shutdown

. /etc/rc.subr

name='routerd_tailscale_mesh'
rcvar="${name}_enable"
command='/usr/local/bin/tailscale'
command_args=" 'up' '--hostname=routerd-mesh' '--auth-key=${TS_AUTHKEY}' '--accept-routes=true' '--accept-dns=false'"
start_precmd="${name}_prestart"

routerd_tailscale_mesh_prestart() {
  mkdir -p '/var/run/routerd'
  'service' 'tailscaled' 'onestart'
}

load_rc_config $name
: ${routerd_tailscale_mesh_enable:="YES"}
run_rc_command "$1"
