# tests/golden/render/freebsd/tailscale-exit-subnet.golden
### rc.conf.d-routerd
# Generated by routerd. Do not edit by hand.
gateway_enable="YES"
ipv6_gateway_enable="YES"
tailscaled_enable="YES"
routerd_tailscale_home_enable="YES"
### rc.d-routerd_tailscale_home
#!/bin/sh
#
# Managed by routerd.
#
# PROVIDE: routerd_tailscale_home
# REQUIRE: NETWORKING tailscaled
# KEYWORD: shutdown

. /etc/rc.subr

name='routerd_tailscale_home'
rcvar="${name}_enable"
command='/usr/local/bin/tailscale'
command_args=" 'up' '--hostname=edge-router' '--auth-key=${TS_AUTHKEY}' '--advertise-exit-node' '--advertise-routes=172.18.0.0/16,192.168.20.0/24' '--accept-dns=false'"
start_precmd="${name}_prestart"

routerd_tailscale_home_prestart() {
  mkdir -p '/var/run/routerd'
  'service' 'tailscaled' 'onestart'
}

load_rc_config $name
: ${routerd_tailscale_home_enable:="YES"}
run_rc_command "$1"
