# tests/golden/render/freebsd/vrrp-tuning-presets.golden
### rc.conf.d-routerd
# Generated by routerd. Do not edit by hand.
gateway_enable="YES"
ipv6_gateway_enable="YES"
routerd_carp_enable="YES"
### rc.d-routerd_carp
#!/bin/sh
#
# Managed by routerd.
#
# PROVIDE: routerd_carp
# REQUIRE: NETWORKING
# KEYWORD: shutdown

. /etc/rc.subr

name="routerd_carp"
rcvar="routerd_carp_enable"
start_cmd="routerd_carp_start"
stop_cmd="routerd_carp_stop"
status_cmd="routerd_carp_status"

routerd_carp_start() {
  kldload carp >/dev/null 2>&1 || true
  sysctl net.inet.carp.preempt='1' >/dev/null
  ifconfig 'ens19' 'inet' 'vhid' '60' 'advbase' '1' 'advskew' '104' 'alias' '192.0.2.10/32'
  ifconfig 'ens19' 'inet' 'vhid' '61' 'advbase' '3' 'advskew' '134' 'alias' '192.0.2.1/32'
}

routerd_carp_stop() {
  ifconfig 'ens19' inet '192.0.2.10/32' -alias >/dev/null 2>&1 || true
  ifconfig 'ens19' inet '192.0.2.1/32' -alias >/dev/null 2>&1 || true
}

routerd_carp_status() {
  true && ifconfig 'ens19' | grep -q 'vhid 60' && ifconfig 'ens19' | grep -q 'vhid 61'
}

load_rc_config $name
: ${routerd_carp_enable:="YES"}
run_rc_command "$1"
