# tests/golden/render/freebsd/wireguard-hub-spoke.golden
### rc.conf.d-routerd
# Generated by routerd. Do not edit by hand.
gateway_enable="YES"
ipv6_gateway_enable="YES"
routerd_wireguard_wg_hub_enable="YES"
### rc.d-routerd_wireguard_wg_hub
#!/bin/sh
#
# Managed by routerd.
#
# PROVIDE: routerd_wireguard_wg_hub
# REQUIRE: NETWORKING
# KEYWORD: shutdown

. /etc/rc.subr

name='routerd_wireguard_wg_hub'
rcvar="${name}_enable"
start_cmd="${name}_start"
stop_cmd="${name}_stop"
status_cmd="${name}_status"

routerd_wireguard_wg_hub_start() {
  kldload if_wg >/dev/null 2>&1 || true
  if ! ifconfig 'wg-hub' >/dev/null 2>&1; then
    ifconfig 'wg-hub' create >/dev/null 2>&1 || ifconfig wg create name 'wg-hub'
  fi
  ifconfig 'wg-hub' mtu '1420'
  wg set 'wg-hub' listen-port '51820' private-key '/usr/local/etc/routerd/secrets/wg-hub.key'
  wg set 'wg-hub' peer 'REPLACE_WITH_SPOKE_A_PUBLIC_KEY' allowed-ips '10.44.0.11/32,172.30.11.0/24' endpoint 'spoke-a.example.net:51820' persistent-keepalive '25'
  wg set 'wg-hub' peer 'REPLACE_WITH_SPOKE_B_PUBLIC_KEY' allowed-ips '10.44.0.12/32,172.30.12.0/24' endpoint 'spoke-b.example.net:51820' persistent-keepalive '25'
  ifconfig 'wg-hub' inet '10.44.0.1/24' alias >/dev/null 2>&1 || true
  ifconfig 'wg-hub' up
}

routerd_wireguard_wg_hub_stop() {
  ifconfig 'wg-hub' destroy >/dev/null 2>&1 || true
}

routerd_wireguard_wg_hub_status() {
  ifconfig 'wg-hub' >/dev/null 2>&1 && wg show 'wg-hub' >/dev/null 2>&1
}

load_rc_config $name
: ${routerd_wireguard_wg_hub_enable:="YES"}
run_rc_command "$1"
