# tests/golden/render/alpine/dhcp-lease-sync-ha.golden
### dnsmasq.conf
# Generated by routerd. Do not edit by hand.
dhcp-leasefile=/var/lib/routerd/dnsmasq/dnsmasq.leases
port=0
no-resolv
bind-dynamic
interface=ens19

interface=ens19
dhcp-range=set:lan-dhcpv4,192.168.30.100,192.168.30.199,12h
dhcp-option=tag:lan-dhcpv4,option:router,192.168.30.254
dhcp-option=tag:lan-dhcpv4,option:dns-server,192.168.30.254

interface=ens19
enable-ra
dhcp-range=set:lan-dhcpv6,fd00:30::100,fd00:30::1ff,constructor:ens19,64,12h
### keepalived.conf
# Generated by routerd. Do not edit by hand.
global_defs {
  router_id routerd
}

vrrp_instance lan_vip {
  state BACKUP
  interface ens19
  virtual_router_id 30
  priority 150
  advert_int 1
  nopreempt
  unicast_peer {
    192.168.30.3
  }
  virtual_ipaddress {
    192.168.30.254/32 dev ens19
  }
}

### openrc-keepalived
#!/sbin/openrc-run
# Generated by routerd. Do not edit by hand.

name='keepalived'
description='routerd managed keepalived VRRP service'
command='/usr/sbin/keepalived'
command_args="'--dont-fork' '--log-console' '--use-file' '/etc/keepalived/keepalived.conf'"
command_background="yes"
pidfile="/run/routerd/openrc/${RC_SVCNAME}.pid"

depend() {
	use net
}

start_pre() {
	checkpath -d -m 0755 /run/routerd/openrc
	checkpath -d -m 0755 '/run/routerd'
	checkpath -d -m 0755 '/var/lib/routerd'
	'/usr/sbin/keepalived' '--config-test' '--use-file' '/etc/keepalived/keepalived.conf'
}

reload() {
	ebegin "Reloading keepalived"
	start-stop-daemon --signal HUP --pidfile "${pidfile}"
	eend $?
}
### openrc-routerd
#!/sbin/openrc-run
# Generated by routerd. Do not edit by hand.

name='routerd'
description='routerd daemon'
command='/usr/local/sbin/routerd'
command_args="'serve' '--config' '/usr/local/etc/routerd/router.yaml' '--socket' '/run/routerd/routerd.sock' '--status-socket' '/run/routerd/routerd-status.sock' '--apply-interval' '60s'"
command_background="yes"
pidfile="/run/routerd/routerd.pid"

depend() {
	use net
}

start_pre() {
	checkpath -d -m 0755 /run/routerd/openrc
	checkpath -d -m 0755 '/run/routerd'
	checkpath -d -m 0755 '/run/routerd/bgp'
	checkpath -d -m 0755 '/run/routerd/dhcpv4-client'
	checkpath -d -m 0755 '/run/routerd/dhcpv6-client'
	checkpath -d -m 0755 '/run/routerd/dns-resolver'
	checkpath -d -m 0755 '/run/routerd/pppoe-client'
	checkpath -d -m 0755 '/var/lib/routerd'
	checkpath -d -m 0755 '/var/log/routerd'
}
### openrc-routerd_dnsmasq
#!/sbin/openrc-run
# Generated by routerd. Do not edit by hand.

name='routerd_dnsmasq'
description='routerd managed dnsmasq DHCP service'
command='/usr/sbin/dnsmasq'
command_args="'--keep-in-foreground' '--user=root' '--group=root' '--conf-file=/usr/local/etc/routerd/dnsmasq.conf' '--pid-file=/run/routerd/dnsmasq.pid'"
command_background="yes"
pidfile="/run/routerd/openrc/${RC_SVCNAME}.pid"

depend() {
	use net
	after routerd
}

start_pre() {
	checkpath -d -m 0755 /run/routerd/openrc
	checkpath -d -m 0755 '/run/routerd'
	checkpath -d -m 0755 '/var/lib/routerd'
}
