# tests/golden/render/alpine/dslite-lan-range-snat.golden
### nftables-nat.nft
# Generated by routerd. Do not edit by hand.
add table ip routerd_nat
flush table ip routerd_nat
table ip routerd_nat {
  comment "routerd.owner=routerd routerd.generation=1"
  chain postrouting {
    type nat hook postrouting priority srcnat; policy accept;
    oifname "ds-lite-lan-range" ip saddr 192.168.160.0/24 ip daddr !=192.168.0.0/16 ip daddr !=172.16.0.0/12 ip daddr !=10.0.0.0/8 snat to 192.168.160.250
  }
}
### openrc-routerd
#!/sbin/openrc-run
# Generated by routerd. Do not edit by hand.

name='routerd'
description='routerd daemon'
command='/usr/local/sbin/routerd'
command_args="'serve'"
command_background="yes"
pidfile="/run/routerd/openrc/${RC_SVCNAME}.pid"

depend() {
	use net
}

start_pre() {
	checkpath -d -m 0755 /run/routerd/openrc
	checkpath -d -m 0755 '/run/routerd'
	checkpath -d -m 0755 '/run/routerd/bgp'
	checkpath -d -m 0755 '/run/routerd/dhcpv4-client'
	checkpath -d -m 0755 '/run/routerd/dhcpv6-client'
	checkpath -d -m 0755 '/run/routerd/dns-resolver'
	checkpath -d -m 0755 '/run/routerd/pppoe-client'
	checkpath -d -m 0755 '/var/lib/routerd'
	checkpath -d -m 0755 '/var/log/routerd'
	'/usr/local/sbin/routerd' 'check'
}
