# tests/golden/render/alpine/example-lan-dns-dhcp.golden
### dnsmasq.conf
# Generated by routerd. Do not edit by hand.
dhcp-leasefile=/var/lib/routerd/dnsmasq/dnsmasq.leases
port=0
no-resolv
bind-dynamic
interface=ens19

interface=ens19
dhcp-range=set:lan-dhcpv4,192.168.30.100,192.168.30.199,12h
dhcp-option=tag:lan-dhcpv4,option:dns-server,192.168.30.1
dhcp-option=tag:lan-dhcpv4,option:domain-name,home.example
dhcp-option=tag:lan-dhcpv4,option:domain-search,home.example
dhcp-host=02:00:00:00:30:10,set:nas,nas,192.168.30.10
### openrc-routerd
#!/sbin/openrc-run
# Generated by routerd. Do not edit by hand.

name='routerd'
description='routerd daemon'
command='/usr/local/sbin/routerd'
command_args="'serve' '--config' '/usr/local/etc/routerd/router.yaml' '--socket' '/run/routerd/routerd.sock' '--status-socket' '/run/routerd/routerd-status.sock' '--apply-interval' '60s'"
command_background="yes"
pidfile="/run/routerd/routerd.pid"

depend() {
	use net
}

start_pre() {
	routerd_stop_client_daemons
	checkpath -d -m 0755 /run/routerd/openrc
	checkpath -d -m 0755 '/run/routerd'
	checkpath -d -m 0755 '/run/routerd/bgp'
	checkpath -d -m 0755 '/run/routerd/dhcpv4-client'
	checkpath -d -m 0755 '/run/routerd/dhcpv6-client'
	checkpath -d -m 0755 '/run/routerd/dns-resolver'
	checkpath -d -m 0755 '/run/routerd/pppoe-client'
	checkpath -d -m 0755 '/var/lib/routerd'
	checkpath -d -m 0755 '/var/log/routerd'
}

stop_post() {
	routerd_stop_client_daemons
}

routerd_stop_client_daemons() {
	pids=''
	for daemon in routerd-dns-resolver routerd-dhcpv4-client routerd-dhcpv6-client routerd-pppoe-client routerd-arp-observer; do
		for pid in $(pidof "${daemon}" 2>/dev/null || true); do
			cmdline=$(tr '\000' ' ' < "/proc/${pid}/cmdline" 2>/dev/null || true)
			case " ${cmdline} " in
				*' --socket /run/routerd/'*|*' --socket /var/run/routerd/'*) pids="${pids} ${pid}" ;;
			esac
		done
	done
	[ -n "${pids}" ] || return 0
	# shellcheck disable=SC2086
	kill -TERM ${pids} 2>/dev/null || true
	i=0
	while [ "${i}" -lt 20 ]; do
		alive=''
		for pid in ${pids}; do
			[ -d "/proc/${pid}" ] && alive="${alive} ${pid}"
		done
		[ -z "${alive}" ] && return 0
		sleep 0.5
		i=$((i + 1))
	done
	# shellcheck disable=SC2086
	kill -KILL ${alive} 2>/dev/null || true
}
### openrc-routerd_dnsmasq
#!/sbin/openrc-run
# Generated by routerd. Do not edit by hand.

name='routerd_dnsmasq'
description='routerd managed dnsmasq DHCP service'
command='/usr/sbin/dnsmasq'
command_args="'--keep-in-foreground' '--user=root' '--group=root' '--conf-file=/usr/local/etc/routerd/dnsmasq.conf' '--pid-file=/run/routerd/dnsmasq.pid'"
command_background="yes"
pidfile="/run/routerd/openrc/${RC_SVCNAME}.pid"

depend() {
	use net
	after routerd
}

start_pre() {
	checkpath -d -m 0755 /run/routerd/openrc
	checkpath -d -m 0755 '/run/routerd'
	checkpath -d -m 0755 '/var/lib/routerd'
}
