# tests/golden/render/alpine/example-lan-dns-dhcp.golden
### dnsmasq.conf
# Generated by routerd. Do not edit by hand.
dhcp-leasefile=/var/lib/routerd/dnsmasq/dnsmasq.leases
port=0
no-resolv
bind-dynamic
interface=ens19

interface=ens19
dhcp-range=set:lan-dhcpv4,192.168.30.100,192.168.30.199,12h
dhcp-option=tag:lan-dhcpv4,option:dns-server,192.168.30.1
dhcp-option=tag:lan-dhcpv4,option:domain-name,home.example
dhcp-option=tag:lan-dhcpv4,option:domain-search,home.example
dhcp-host=02:00:00:00:30:10,set:nas,nas,192.168.30.10
### openrc-routerd
#!/sbin/openrc-run
# Generated by routerd. Do not edit by hand.

name='routerd'
description='routerd daemon'
command='/usr/local/sbin/routerd'
command_args="'serve' '--config' '/usr/local/etc/routerd/router.yaml' '--socket' '/run/routerd/routerd.sock' '--status-socket' '/run/routerd/routerd-status.sock' '--apply-interval' '60s'"
command_background="yes"
pidfile="/run/routerd/routerd.pid"

depend() {
	use net
}

start_pre() {
	checkpath -d -m 0755 /run/routerd/openrc
	checkpath -d -m 0755 '/run/routerd'
	checkpath -d -m 0755 '/run/routerd/bgp'
	checkpath -d -m 0755 '/run/routerd/dhcpv4-client'
	checkpath -d -m 0755 '/run/routerd/dhcpv6-client'
	checkpath -d -m 0755 '/run/routerd/dns-resolver'
	checkpath -d -m 0755 '/run/routerd/pppoe-client'
	checkpath -d -m 0755 '/var/lib/routerd'
	checkpath -d -m 0755 '/var/log/routerd'
}
### openrc-routerd_dnsmasq
#!/sbin/openrc-run
# Generated by routerd. Do not edit by hand.

name='routerd_dnsmasq'
description='routerd managed dnsmasq DHCP service'
command='/usr/sbin/dnsmasq'
command_args="'--keep-in-foreground' '--user=root' '--group=root' '--conf-file=/usr/local/etc/routerd/dnsmasq.conf' '--pid-file=/run/routerd/dnsmasq.pid'"
command_background="yes"
pidfile="/run/routerd/openrc/${RC_SVCNAME}.pid"

depend() {
	use net
	after routerd
}

start_pre() {
	checkpath -d -m 0755 /run/routerd/openrc
	checkpath -d -m 0755 '/run/routerd'
	checkpath -d -m 0755 '/var/lib/routerd'
}
