# tests/golden/render/alpine/freebsd-edge.golden
### dnsmasq.conf
# Generated by routerd. Do not edit by hand.
dhcp-leasefile=/var/lib/routerd/dnsmasq/dnsmasq.leases
port=0
no-resolv
bind-dynamic
interface=vtnet1
except-interface=vtnet0
except-interface=vtnet2

dhcp-authoritative
dhcp-range=set:lan-dhcpv4,192.0.2.110,192.0.2.119,255.255.255.0,12h
dhcp-option=tag:lan-dhcpv4,option:router,192.0.2.1
dhcp-option=tag:lan-dhcpv4,option:dns-server,192.0.2.1
dhcp-host=02:00:00:00:01:50,192.0.2.115,printer,infinite
### openrc-routerd_dhcpv4_client_mgmt_dhcpv4
#!/sbin/openrc-run
# Generated by routerd. Do not edit by hand.

name='routerd_dhcpv4_client_mgmt_dhcpv4'
description='routerd DHCPv4 client mgmt-dhcpv4'
command='/usr/local/sbin/routerd-dhcpv4-client'
command_args="'daemon' '--resource' 'mgmt-dhcpv4' '--interface' 'vtnet2'"
command_background="yes"
pidfile="/run/routerd/openrc/${RC_SVCNAME}.pid"

depend() {
	use net
}

start_pre() {
	checkpath -d -m 0755 /run/routerd/openrc
	checkpath -d -m 0755 '/run/routerd/dhcpv4-client'
	checkpath -d -m 0755 '/var/lib/routerd/dhcpv4-client'
	checkpath -d -m 0755 '/var/log/routerd'
}
### openrc-routerd_dhcpv4_client_wan_dhcpv4
#!/sbin/openrc-run
# Generated by routerd. Do not edit by hand.

name='routerd_dhcpv4_client_wan_dhcpv4'
description='routerd DHCPv4 client wan-dhcpv4'
command='/usr/local/sbin/routerd-dhcpv4-client'
command_args="'daemon' '--resource' 'wan-dhcpv4' '--interface' 'vtnet0'"
command_background="yes"
pidfile="/run/routerd/openrc/${RC_SVCNAME}.pid"

depend() {
	use net
}

start_pre() {
	checkpath -d -m 0755 /run/routerd/openrc
	checkpath -d -m 0755 '/run/routerd/dhcpv4-client'
	checkpath -d -m 0755 '/var/lib/routerd/dhcpv4-client'
	checkpath -d -m 0755 '/var/log/routerd'
}
### openrc-routerd_dhcpv6_client_wan_pd
#!/sbin/openrc-run
# Generated by routerd. Do not edit by hand.

name='routerd_dhcpv6_client_wan_pd'
description='routerd DHCPv6-PD client wan-pd'
command='/usr/local/sbin/routerd-dhcpv6-client'
command_args="'daemon' '--resource' 'wan-pd' '--interface' 'vtnet0' '--socket' '/run/routerd/dhcpv6-client/wan-pd.sock' '--lease-file' '/var/lib/routerd/dhcpv6-client/wan-pd/lease.json' '--event-file' '/var/lib/routerd/dhcpv6-client/wan-pd/events.jsonl'"
command_background="yes"
pidfile="/run/routerd/openrc/${RC_SVCNAME}.pid"

depend() {
	use net
}

start_pre() {
	checkpath -d -m 0755 /run/routerd/openrc
	checkpath -d -m 0755 '/run/routerd/dhcpv6-client'
	checkpath -d -m 0755 '/var/lib/routerd/dhcpv6-client'
	checkpath -d -m 0755 '/var/lib/routerd/dhcpv6-client/wan-pd'
	checkpath -d -m 0755 '/var/log/routerd'
}
### openrc-routerd_dnsmasq
#!/sbin/openrc-run
# Generated by routerd. Do not edit by hand.

name='routerd_dnsmasq'
description='routerd managed dnsmasq DHCP service'
command='/usr/sbin/dnsmasq'
command_args="'--keep-in-foreground' '--user=root' '--group=root' '--conf-file=/usr/local/etc/routerd/dnsmasq.conf' '--pid-file=/run/routerd/dnsmasq.pid'"
command_background="yes"
pidfile="/run/routerd/openrc/${RC_SVCNAME}.pid"

depend() {
	use net
	after routerd
}

start_pre() {
	checkpath -d -m 0755 /run/routerd/openrc
	checkpath -d -m 0755 '/run/routerd'
	checkpath -d -m 0755 '/var/lib/routerd'
}
### openrc-routerd_dpi_classifier
#!/sbin/openrc-run
# Generated by routerd. Do not edit by hand.

name='routerd_dpi_classifier'
description='routerd DPI classifier'
command='/usr/local/sbin/routerd-dpi-classifier'
command_args="'daemon' '--socket' '/var/run/routerd/dpi-classifier/default.sock' '--ndpi-reader' '/usr/local/bin/ndpiReader'"
command_background="yes"
pidfile="/run/routerd/openrc/${RC_SVCNAME}.pid"

depend() {
	use net
}

start_pre() {
	checkpath -d -m 0755 /run/routerd/openrc
	checkpath -d -m 0755 '/run/routerd/dpi-classifier'
}
### openrc-routerd_firewall_logger
#!/sbin/openrc-run
# Generated by routerd. Do not edit by hand.

name='routerd_firewall_logger'
description='routerd firewall logger'
command='/usr/local/sbin/routerd-firewall-logger'
command_args="'daemon' '--path' '/var/db/routerd/firewall-logs.db' '--pflog-interface' 'pflog0' '--dpi-socket' '/var/run/routerd/dpi-classifier/default.sock'"
command_background="yes"
pidfile="/run/routerd/openrc/${RC_SVCNAME}.pid"

depend() {
	use net
	after routerd_dpi_classifier
}

start_pre() {
	checkpath -d -m 0755 /run/routerd/openrc
	checkpath -d -m 0755 '/run/routerd'
	checkpath -d -m 0755 '/var/lib/routerd'
}
### openrc-routerd_tailscale_home
#!/sbin/openrc-run
# Generated by routerd. Do not edit by hand.

name='routerd_tailscale_home'
description='routerd Tailscale node home'
command='/usr/bin/tailscale'
command_args="'up' '--hostname=freebsd-edge' '--advertise-routes=192.0.2.0/24' '--accept-dns=false'"

depend() {
	use net
	after tailscaled
}

start_pre() {
	checkpath -d -m 0755 /run/routerd/openrc
	checkpath -d -m 0755 '/run/routerd'
	'rc-service' 'tailscaled' 'start'
}
