# tests/golden/render/alpine/k8s-routerd-vip-alpine.golden
### keepalived.conf
# Generated by routerd. Do not edit by hand.
global_defs {
  router_id routerd
}

vrrp_instance k8s_api_vip {
  state BACKUP
  interface eth0
  virtual_router_id 50
  priority 150
  advert_int 1
  unicast_peer {
    192.168.1.53
  }
  virtual_ipaddress {
    192.168.1.248/32 dev eth0
  }
}

### nftables-filter.nft.error
firewall.routerd.net/v1alpha1/IngressService/kubernetes-api backend 0 address: must resolve to an IPv4 address
### nftables-nat.nft.error
firewall.routerd.net/v1alpha1/IngressService/kubernetes-api backend 0 address: must resolve to an IPv4 address
### openrc-keepalived
#!/sbin/openrc-run
# Generated by routerd. Do not edit by hand.

name='keepalived'
description='routerd managed keepalived VRRP service'
command='/usr/sbin/keepalived'
command_args="'--dont-fork' '--log-console' '--use-file' '/etc/keepalived/keepalived.conf'"
command_background="yes"
pidfile="/run/routerd/openrc/${RC_SVCNAME}.pid"

depend() {
	use net
}

start_pre() {
	checkpath -d -m 0755 /run/routerd/openrc
	checkpath -d -m 0755 '/run/routerd'
	checkpath -d -m 0755 '/var/lib/routerd'
	'/usr/sbin/keepalived' '--config-test' '--use-file' '/etc/keepalived/keepalived.conf'
}

reload() {
	ebegin "Reloading keepalived"
	start-stop-daemon --signal HUP --pidfile "${pidfile}"
	eend $?
}
### openrc-routerd
#!/sbin/openrc-run
# Generated by routerd. Do not edit by hand.

name='routerd'
description='routerd daemon'
command='/usr/local/sbin/routerd'
command_args="'serve' '--config' '/usr/local/etc/routerd/router.yaml' '--socket' '/run/routerd/routerd.sock' '--status-socket' '/run/routerd/routerd-status.sock' '--apply-interval' '60s'"
command_background="yes"
pidfile="/run/routerd/routerd.pid"

depend() {
	use net
}

start_pre() {
	routerd_stop_client_daemons
	checkpath -d -m 0755 /run/routerd/openrc
	checkpath -d -m 0755 '/run/routerd'
	checkpath -d -m 0755 '/run/routerd/bgp'
	checkpath -d -m 0755 '/run/routerd/dhcpv4-client'
	checkpath -d -m 0755 '/run/routerd/dhcpv6-client'
	checkpath -d -m 0755 '/run/routerd/dns-resolver'
	checkpath -d -m 0755 '/run/routerd/pppoe-client'
	checkpath -d -m 0755 '/var/lib/routerd'
	checkpath -d -m 0755 '/var/log/routerd'
}

stop_post() {
	routerd_stop_client_daemons
}

routerd_stop_client_daemons() {
	pids=''
	for daemon in routerd-dns-resolver routerd-dhcpv4-client routerd-dhcpv6-client routerd-pppoe-client routerd-arp-observer; do
		for pid in $(pidof "${daemon}" 2>/dev/null || true); do
			cmdline=$(tr '\000' ' ' < "/proc/${pid}/cmdline" 2>/dev/null || true)
			case " ${cmdline} " in
				*' --socket /run/routerd/'*|*' --socket /var/run/routerd/'*) pids="${pids} ${pid}" ;;
			esac
		done
	done
	[ -n "${pids}" ] || return 0
	# shellcheck disable=SC2086
	kill -TERM ${pids} 2>/dev/null || true
	i=0
	while [ "${i}" -lt 20 ]; do
		alive=''
		for pid in ${pids}; do
			[ -d "/proc/${pid}" ] && alive="${alive} ${pid}"
		done
		[ -z "${alive}" ] && return 0
		sleep 0.5
		i=$((i + 1))
	done
	# shellcheck disable=SC2086
	kill -KILL ${alive} 2>/dev/null || true
}
