# tests/golden/render/alpine/kubernetes-api-vip.golden
### keepalived.conf
# Generated by routerd. Do not edit by hand.
global_defs {
  router_id routerd
}

vrrp_instance k8s_api_vip {
  state BACKUP
  interface ens20
  virtual_router_id 70
  priority 150
  advert_int 1
  nopreempt
  unicast_peer {
    192.168.70.3
  }
  virtual_ipaddress {
    192.168.70.10/32 dev ens20
  }
}

### nftables-filter.nft.error
firewall.routerd.net/v1alpha1/IngressService/kubernetes-api backend 0 address: must resolve to an IPv4 address
### nftables-nat.nft.error
firewall.routerd.net/v1alpha1/IngressService/kubernetes-api backend 0 address: must resolve to an IPv4 address
### openrc-keepalived
#!/sbin/openrc-run
# Generated by routerd. Do not edit by hand.

name='keepalived'
description='routerd managed keepalived VRRP service'
command='/usr/sbin/keepalived'
command_args="'--dont-fork' '--log-console' '--use-file' '/etc/keepalived/keepalived.conf'"
command_background="yes"
pidfile="/run/routerd/openrc/${RC_SVCNAME}.pid"

depend() {
	use net
}

start_pre() {
	checkpath -d -m 0755 /run/routerd/openrc
	checkpath -d -m 0755 '/run/routerd'
	checkpath -d -m 0755 '/var/lib/routerd'
	'/usr/sbin/keepalived' '--config-test' '--use-file' '/etc/keepalived/keepalived.conf'
}

reload() {
	ebegin "Reloading keepalived"
	start-stop-daemon --signal HUP --pidfile "${pidfile}"
	eend $?
}
### openrc-routerd
#!/sbin/openrc-run
# Generated by routerd. Do not edit by hand.

name='routerd'
description='routerd daemon'
command='/usr/local/sbin/routerd'
command_args="'serve'"
command_background="yes"
pidfile="/run/routerd/openrc/${RC_SVCNAME}.pid"

depend() {
	use net
}

start_pre() {
	checkpath -d -m 0755 /run/routerd/openrc
	checkpath -d -m 0755 '/run/routerd'
	checkpath -d -m 0755 '/run/routerd/bgp'
	checkpath -d -m 0755 '/run/routerd/dhcpv4-client'
	checkpath -d -m 0755 '/run/routerd/dhcpv6-client'
	checkpath -d -m 0755 '/run/routerd/dns-resolver'
	checkpath -d -m 0755 '/run/routerd/pppoe-client'
	checkpath -d -m 0755 '/var/lib/routerd'
	checkpath -d -m 0755 '/var/log/routerd'
	'/usr/local/sbin/routerd' 'check'
}
### openrc-routerd_dns_resolver_lan
#!/sbin/openrc-run
# Generated by routerd. Do not edit by hand.

name='routerd_dns_resolver_lan'
description='routerd DNS resolver lan'
command='/usr/local/sbin/routerd-dns-resolver'
command_args="'daemon' '--resource' 'lan' '--config-file' '/var/lib/routerd/dns-resolver/lan/config.json'"
command_background="yes"
pidfile="/run/routerd/openrc/${RC_SVCNAME}.pid"

depend() {
	use net
}

start_pre() {
	checkpath -d -m 0755 /run/routerd/openrc
	checkpath -d -m 0755 '/run/routerd/dns-resolver'
	checkpath -d -m 0755 '/var/lib/routerd/dns-resolver'
	checkpath -d -m 0755 '/var/lib/routerd/dns-resolver/lan'
	checkpath -d -m 0755 '/var/log/routerd'
}
