# tests/golden/render/alpine/kubernetes-api-vip-external-dns.golden
### nftables-nat.nft
# Generated by routerd. Do not edit by hand.
add table ip routerd_nat
flush table ip routerd_nat
table ip routerd_nat {
  comment "routerd.owner=routerd routerd.generation=1"
  chain prerouting {
    type nat hook prerouting priority dstnat; policy accept;
    iifname "eth0" ip daddr 192.0.2.10 tcp dport 6443 counter dnat to 192.0.2.11:6443 comment "routerd IngressService kubernetes-api"
  }
  chain postrouting {
    type nat hook postrouting priority srcnat; policy accept;
  }
}
### openrc-routerd
#!/sbin/openrc-run
# Generated by routerd. Do not edit by hand.

name='routerd'
description='routerd daemon'
command='/usr/local/sbin/routerd'
command_args="'serve' '--config' '/usr/local/etc/routerd/router.yaml' '--socket' '/run/routerd/routerd.sock' '--status-socket' '/run/routerd/routerd-status.sock' '--apply-interval' '60s'"
command_background="yes"
pidfile="/run/routerd/routerd.pid"

depend() {
	use net
}

start_pre() {
	routerd_stop_client_daemons
	checkpath -d -m 0755 /run/routerd/openrc
	checkpath -d -m 0755 '/run/routerd'
	checkpath -d -m 0755 '/run/routerd/bgp'
	checkpath -d -m 0755 '/run/routerd/dhcpv4-client'
	checkpath -d -m 0755 '/run/routerd/dhcpv6-client'
	checkpath -d -m 0755 '/run/routerd/dns-resolver'
	checkpath -d -m 0755 '/run/routerd/pppoe-client'
	checkpath -d -m 0755 '/var/lib/routerd'
	checkpath -d -m 0755 '/var/log/routerd'
}

stop_post() {
	routerd_stop_client_daemons
}

routerd_stop_client_daemons() {
	pids=''
	for daemon in routerd-dns-resolver routerd-dhcpv4-client routerd-dhcpv6-client routerd-pppoe-client routerd-arp-observer; do
		for pid in $(pidof "${daemon}" 2>/dev/null || true); do
			cmdline=$(tr '\000' ' ' < "/proc/${pid}/cmdline" 2>/dev/null || true)
			case " ${cmdline} " in
				*' --socket /run/routerd/'*|*' --socket /var/run/routerd/'*) pids="${pids} ${pid}" ;;
			esac
		done
	done
	[ -n "${pids}" ] || return 0
	# shellcheck disable=SC2086
	kill -TERM ${pids} 2>/dev/null || true
	i=0
	while [ "${i}" -lt 20 ]; do
		alive=''
		for pid in ${pids}; do
			[ -d "/proc/${pid}" ] && alive="${alive} ${pid}"
		done
		[ -z "${alive}" ] && return 0
		sleep 0.5
		i=$((i + 1))
	done
	# shellcheck disable=SC2086
	kill -KILL ${alive} 2>/dev/null || true
}
