# tests/golden/render/alpine/multi-wan-home.golden
### nftables-nat.nft
# Generated by routerd. Do not edit by hand.
add table ip routerd_nat
flush table ip routerd_nat
table ip routerd_nat {
  chain postrouting {
    type nat hook postrouting priority srcnat; policy accept;
    # net.routerd.net/v1alpha1/NAT44Rule/lan-to-selected-wan skipped: egressInterface is resolved by ipv4-default at runtime
  }
}
### openrc-routerd
#!/sbin/openrc-run
# Generated by routerd. Do not edit by hand.

name='routerd'
description='routerd daemon'
command='/usr/local/sbin/routerd'
command_args="'serve'"
command_background="yes"
pidfile="/run/routerd/openrc/${RC_SVCNAME}.pid"

depend() {
	use net
}

start_pre() {
	checkpath -d -m 0755 /run/routerd/openrc
	checkpath -d -m 0755 '/run/routerd'
	checkpath -d -m 0755 '/run/routerd/dhcpv4-client'
	checkpath -d -m 0755 '/run/routerd/dhcpv6-client'
	checkpath -d -m 0755 '/run/routerd/dns-resolver'
	checkpath -d -m 0755 '/run/routerd/pppoe-client'
	checkpath -d -m 0755 '/var/lib/routerd'
	checkpath -d -m 0755 '/var/log/routerd'
	'/usr/local/sbin/routerd' 'check'
}
### openrc-routerd_healthcheck_internet_via_dslite_a
#!/sbin/openrc-run
# Generated by routerd. Do not edit by hand.

name='routerd_healthcheck_internet_via_dslite_a'
description='routerd healthcheck internet-via-dslite-a'
command='/usr/local/sbin/routerd-healthcheck'
command_args="'daemon' '--resource' 'internet-via-dslite-a' '--target' '1.1.1.1' '--protocol' 'tcp' '--address-family' 'ipv4' '--port' '443' '--interval' '30s' '--timeout' '3s' '--healthy-threshold' '2' '--unhealthy-threshold' '2' '--socket' '/run/routerd/healthcheck/internet-via-dslite-a.sock' '--state-file' '/var/lib/routerd/healthcheck/internet-via-dslite-a/state.json' '--event-file' '/var/log/routerd/healthcheck/internet-via-dslite-a/events.jsonl'"
command_background="yes"
pidfile="/run/routerd/openrc/${RC_SVCNAME}.pid"

depend() {
	use net
}

start_pre() {
	checkpath -d -m 0755 /run/routerd/openrc
	checkpath -d -m 0755 '/run/routerd/healthcheck'
	checkpath -d -m 0755 '/var/lib/routerd/healthcheck'
	checkpath -d -m 0755 '/var/log/routerd'
}
### openrc-routerd_healthcheck_internet_via_dslite_b
#!/sbin/openrc-run
# Generated by routerd. Do not edit by hand.

name='routerd_healthcheck_internet_via_dslite_b'
description='routerd healthcheck internet-via-dslite-b'
command='/usr/local/sbin/routerd-healthcheck'
command_args="'daemon' '--resource' 'internet-via-dslite-b' '--target' '1.1.1.1' '--protocol' 'tcp' '--address-family' 'ipv4' '--port' '443' '--interval' '30s' '--timeout' '3s' '--healthy-threshold' '2' '--unhealthy-threshold' '2' '--socket' '/run/routerd/healthcheck/internet-via-dslite-b.sock' '--state-file' '/var/lib/routerd/healthcheck/internet-via-dslite-b/state.json' '--event-file' '/var/log/routerd/healthcheck/internet-via-dslite-b/events.jsonl'"
command_background="yes"
pidfile="/run/routerd/openrc/${RC_SVCNAME}.pid"

depend() {
	use net
}

start_pre() {
	checkpath -d -m 0755 /run/routerd/openrc
	checkpath -d -m 0755 '/run/routerd/healthcheck'
	checkpath -d -m 0755 '/var/lib/routerd/healthcheck'
	checkpath -d -m 0755 '/var/log/routerd'
}
### openrc-routerd_healthcheck_internet_via_dslite_c
#!/sbin/openrc-run
# Generated by routerd. Do not edit by hand.

name='routerd_healthcheck_internet_via_dslite_c'
description='routerd healthcheck internet-via-dslite-c'
command='/usr/local/sbin/routerd-healthcheck'
command_args="'daemon' '--resource' 'internet-via-dslite-c' '--target' '1.1.1.1' '--protocol' 'tcp' '--address-family' 'ipv4' '--port' '443' '--interval' '30s' '--timeout' '3s' '--healthy-threshold' '2' '--unhealthy-threshold' '2' '--socket' '/run/routerd/healthcheck/internet-via-dslite-c.sock' '--state-file' '/var/lib/routerd/healthcheck/internet-via-dslite-c/state.json' '--event-file' '/var/log/routerd/healthcheck/internet-via-dslite-c/events.jsonl'"
command_background="yes"
pidfile="/run/routerd/openrc/${RC_SVCNAME}.pid"

depend() {
	use net
}

start_pre() {
	checkpath -d -m 0755 /run/routerd/openrc
	checkpath -d -m 0755 '/run/routerd/healthcheck'
	checkpath -d -m 0755 '/var/lib/routerd/healthcheck'
	checkpath -d -m 0755 '/var/log/routerd'
}
### openrc-routerd_healthcheck_internet_via_dslite_ra
#!/sbin/openrc-run
# Generated by routerd. Do not edit by hand.

name='routerd_healthcheck_internet_via_dslite_ra'
description='routerd healthcheck internet-via-dslite-ra'
command='/usr/local/sbin/routerd-healthcheck'
command_args="'daemon' '--resource' 'internet-via-dslite-ra' '--target' '1.1.1.1' '--protocol' 'tcp' '--address-family' 'ipv4' '--port' '443' '--interval' '30s' '--timeout' '3s' '--healthy-threshold' '2' '--unhealthy-threshold' '2' '--socket' '/run/routerd/healthcheck/internet-via-dslite-ra.sock' '--state-file' '/var/lib/routerd/healthcheck/internet-via-dslite-ra/state.json' '--event-file' '/var/log/routerd/healthcheck/internet-via-dslite-ra/events.jsonl'"
command_background="yes"
pidfile="/run/routerd/openrc/${RC_SVCNAME}.pid"

depend() {
	use net
}

start_pre() {
	checkpath -d -m 0755 /run/routerd/openrc
	checkpath -d -m 0755 '/run/routerd/healthcheck'
	checkpath -d -m 0755 '/var/lib/routerd/healthcheck'
	checkpath -d -m 0755 '/var/log/routerd'
}
### openrc-routerd_healthcheck_internet_via_hgw_direct
#!/sbin/openrc-run
# Generated by routerd. Do not edit by hand.

name='routerd_healthcheck_internet_via_hgw_direct'
description='routerd healthcheck internet-via-hgw-direct'
command='/usr/local/sbin/routerd-healthcheck'
command_args="'daemon' '--resource' 'internet-via-hgw-direct' '--target' '1.1.1.1' '--protocol' 'tcp' '--address-family' 'ipv4' '--port' '443' '--interval' '30s' '--timeout' '3s' '--healthy-threshold' '2' '--unhealthy-threshold' '2' '--socket' '/run/routerd/healthcheck/internet-via-hgw-direct.sock' '--state-file' '/var/lib/routerd/healthcheck/internet-via-hgw-direct/state.json' '--event-file' '/var/log/routerd/healthcheck/internet-via-hgw-direct/events.jsonl'"
command_background="yes"
pidfile="/run/routerd/openrc/${RC_SVCNAME}.pid"

depend() {
	use net
}

start_pre() {
	checkpath -d -m 0755 /run/routerd/openrc
	checkpath -d -m 0755 '/run/routerd/healthcheck'
	checkpath -d -m 0755 '/var/lib/routerd/healthcheck'
	checkpath -d -m 0755 '/var/log/routerd'
}
### openrc-routerd_healthcheck_internet_via_pppoe
#!/sbin/openrc-run
# Generated by routerd. Do not edit by hand.

name='routerd_healthcheck_internet_via_pppoe'
description='routerd healthcheck internet-via-pppoe'
command='/usr/local/sbin/routerd-healthcheck'
command_args="'daemon' '--resource' 'internet-via-pppoe' '--target' '1.1.1.1' '--protocol' 'tcp' '--address-family' 'ipv4' '--source-interface' 'ppp-flets' '--port' '443' '--interval' '30s' '--timeout' '3s' '--healthy-threshold' '2' '--unhealthy-threshold' '2' '--socket' '/run/routerd/healthcheck/internet-via-pppoe.sock' '--state-file' '/var/lib/routerd/healthcheck/internet-via-pppoe/state.json' '--event-file' '/var/log/routerd/healthcheck/internet-via-pppoe/events.jsonl'"
command_background="yes"
pidfile="/run/routerd/openrc/${RC_SVCNAME}.pid"

depend() {
	use net
}

start_pre() {
	checkpath -d -m 0755 /run/routerd/openrc
	checkpath -d -m 0755 '/run/routerd/healthcheck'
	checkpath -d -m 0755 '/var/lib/routerd/healthcheck'
	checkpath -d -m 0755 '/var/log/routerd'
}
### openrc-routerd_pppoe_client_pppoe_flets
#!/sbin/openrc-run
# Generated by routerd. Do not edit by hand.

name='routerd_pppoe_client_pppoe_flets'
description='routerd PPPoE client pppoe-flets'
command='/usr/local/sbin/routerd-pppoe-client'
command_args="'daemon' '--resource' 'pppoe-flets' '--interface' 'ens18' '--username' 'open@open.ad.jp' '--password-file' '/usr/local/etc/routerd/secrets/pppoe-flets.password' '--mtu' '1454' '--mru' '1454'"
command_background="yes"
pidfile="/run/routerd/openrc/${RC_SVCNAME}.pid"

depend() {
	use net
}

start_pre() {
	checkpath -d -m 0755 /run/routerd/openrc
	checkpath -d -m 0755 '/run/routerd/pppoe-client'
	checkpath -d -m 0755 '/var/lib/routerd/pppoe-client'
	checkpath -d -m 0755 '/var/log/routerd'
}
