# tests/golden/render/alpine/multi-wan-home.golden
### nftables-nat.nft
# Generated by routerd. Do not edit by hand.
add table ip routerd_nat
flush table ip routerd_nat
table ip routerd_nat {
  chain postrouting {
    type nat hook postrouting priority srcnat; policy accept;
    # net.routerd.net/v1alpha1/NAT44Rule/lan-to-selected-wan skipped: egressInterface is resolved by ipv4-default at runtime
  }
}
### openrc-routerd
#!/sbin/openrc-run
# Generated by routerd. Do not edit by hand.

name='routerd'
description='routerd daemon'
command='/usr/local/sbin/routerd'
command_args="'serve'"
command_background="yes"
pidfile="/run/routerd/openrc/${RC_SVCNAME}.pid"

depend() {
	use net
}

start_pre() {
	checkpath -d -m 0755 /run/routerd/openrc
	checkpath -d -m 0755 '/run/routerd'
	checkpath -d -m 0755 '/run/routerd/dhcpv4-client'
	checkpath -d -m 0755 '/run/routerd/dhcpv6-client'
	checkpath -d -m 0755 '/run/routerd/dns-resolver'
	checkpath -d -m 0755 '/run/routerd/pppoe-client'
	checkpath -d -m 0755 '/var/lib/routerd'
	checkpath -d -m 0755 '/var/log/routerd'
	'/usr/local/sbin/routerd' 'check'
}
