# tests/golden/render/alpine/router-lab.golden
### dnsmasq.conf
# Generated by routerd. Do not edit by hand.
dhcp-leasefile=/var/lib/routerd/dnsmasq/dnsmasq.leases
port=0
no-resolv
bind-dynamic
interface=ens19
except-interface=ens18
except-interface=ens20

dhcp-authoritative
dhcp-range=set:lan-dhcpv4,192.168.10.100,192.168.10.199,255.255.255.0,12h
dhcp-option=tag:lan-dhcpv4,option:router,192.168.10.3
dhcp-option=tag:lan-dhcpv4,option:dns-server,192.168.10.3
dhcp-host=02:00:00:00:01:50,192.168.10.150,printer,infinite
### openrc-routerd_dhcpv4_client_mgmt_dhcpv4
#!/sbin/openrc-run
# Generated by routerd. Do not edit by hand.

name='routerd_dhcpv4_client_mgmt_dhcpv4'
description='routerd DHCPv4 client mgmt-dhcpv4'
command='/usr/local/sbin/routerd-dhcpv4-client'
command_args="'daemon' '--resource' 'mgmt-dhcpv4' '--interface' 'ens20'"
command_background="yes"
pidfile="/run/routerd/openrc/${RC_SVCNAME}.pid"

depend() {
	use net
}

start_pre() {
	checkpath -d -m 0755 /run/routerd/openrc
	checkpath -d -m 0755 '/run/routerd/dhcpv4-client'
	checkpath -d -m 0755 '/var/lib/routerd/dhcpv4-client'
	checkpath -d -m 0755 '/var/log/routerd'
}
### openrc-routerd_dhcpv4_client_wan_dhcpv4
#!/sbin/openrc-run
# Generated by routerd. Do not edit by hand.

name='routerd_dhcpv4_client_wan_dhcpv4'
description='routerd DHCPv4 client wan-dhcpv4'
command='/usr/local/sbin/routerd-dhcpv4-client'
command_args="'daemon' '--resource' 'wan-dhcpv4' '--interface' 'ens18'"
command_background="yes"
pidfile="/run/routerd/openrc/${RC_SVCNAME}.pid"

depend() {
	use net
}

start_pre() {
	checkpath -d -m 0755 /run/routerd/openrc
	checkpath -d -m 0755 '/run/routerd/dhcpv4-client'
	checkpath -d -m 0755 '/var/lib/routerd/dhcpv4-client'
	checkpath -d -m 0755 '/var/log/routerd'
}
### openrc-routerd_dhcpv6_client_wan_pd
#!/sbin/openrc-run
# Generated by routerd. Do not edit by hand.

name='routerd_dhcpv6_client_wan_pd'
description='routerd DHCPv6-PD client wan-pd'
command='/usr/local/sbin/routerd-dhcpv6-client'
command_args="'daemon' '--resource' 'wan-pd' '--interface' 'ens18' '--socket' '/run/routerd/dhcpv6-client/wan-pd.sock' '--lease-file' '/var/lib/routerd/dhcpv6-client/wan-pd/lease.json' '--event-file' '/var/lib/routerd/dhcpv6-client/wan-pd/events.jsonl'"
command_background="yes"
pidfile="/run/routerd/openrc/${RC_SVCNAME}.pid"

depend() {
	use net
}

start_pre() {
	checkpath -d -m 0755 /run/routerd/openrc
	checkpath -d -m 0755 '/run/routerd/dhcpv6-client'
	checkpath -d -m 0755 '/var/lib/routerd/dhcpv6-client'
	checkpath -d -m 0755 '/var/lib/routerd/dhcpv6-client/wan-pd'
	checkpath -d -m 0755 '/var/log/routerd'
}
### openrc-routerd_dns_resolver_lan_resolver
#!/sbin/openrc-run
# Generated by routerd. Do not edit by hand.

name='routerd_dns_resolver_lan_resolver'
description='routerd DNS resolver lan-resolver'
command='/usr/local/sbin/routerd-dns-resolver'
command_args="'daemon' '--resource' 'lan-resolver' '--config-file' '/var/lib/routerd/dns-resolver/lan-resolver/config.json'"
command_background="yes"
pidfile="/run/routerd/openrc/${RC_SVCNAME}.pid"

depend() {
	use net
}

start_pre() {
	checkpath -d -m 0755 /run/routerd/openrc
	checkpath -d -m 0755 '/run/routerd/dns-resolver'
	checkpath -d -m 0755 '/var/lib/routerd/dns-resolver'
	checkpath -d -m 0755 '/var/lib/routerd/dns-resolver/lan-resolver'
	checkpath -d -m 0755 '/var/log/routerd'
}
### openrc-routerd_dnsmasq
#!/sbin/openrc-run
# Generated by routerd. Do not edit by hand.

name='routerd_dnsmasq'
description='routerd managed dnsmasq DHCP service'
command='/usr/sbin/dnsmasq'
command_args="'--keep-in-foreground' '--user=root' '--group=root' '--conf-file=/usr/local/etc/routerd/dnsmasq.conf' '--pid-file=/run/routerd/dnsmasq.pid'"
command_background="yes"
pidfile="/run/routerd/openrc/${RC_SVCNAME}.pid"

depend() {
	use net
	after routerd
}

start_pre() {
	checkpath -d -m 0755 /run/routerd/openrc
	checkpath -d -m 0755 '/run/routerd'
	checkpath -d -m 0755 '/var/lib/routerd'
}
