# tests/golden/render/alpine/tailscale-exit-subnet.golden
### openrc-routerd_tailscale_home
#!/sbin/openrc-run
# Generated by routerd. Do not edit by hand.

name='routerd_tailscale_home'
description='routerd Tailscale node home'
command='/usr/bin/tailscale'
command_args="'up' '--hostname=edge-router' '--auth-key=${TS_AUTHKEY}' '--advertise-exit-node' '--advertise-routes=172.18.0.0/16,192.168.20.0/24' '--accept-dns=false'"

depend() {
	use net
	after tailscaled
}

start_pre() {
	checkpath -d -m 0755 /run/routerd/openrc
	checkpath -d -m 0755 '/run/routerd'
	'rc-service' 'tailscaled' 'start'
}
