# tests/golden/render/alpine/tailscale-minimal.golden
### openrc-routerd_tailscale_mesh
#!/sbin/openrc-run
# Generated by routerd. Do not edit by hand.

name='routerd_tailscale_mesh'
description='routerd Tailscale node mesh'
command='/usr/bin/tailscale'
command_args="'up' '--hostname=routerd-mesh' '--auth-key=${TS_AUTHKEY}' '--accept-routes=true' '--accept-dns=false'"

depend() {
	use net
	after tailscaled
}

start_pre() {
	checkpath -d -m 0755 /run/routerd/openrc
	checkpath -d -m 0755 '/run/routerd'
	'rc-service' 'tailscaled' 'start'
}
