# tests/golden/render/alpine/vrrp-tuning-presets.golden
### keepalived.conf
# Generated by routerd. Do not edit by hand.
global_defs {
  router_id routerd
}

vrrp_instance home_lan_vip {
  state BACKUP
  interface ens19
  virtual_router_id 61
  priority 120
  advert_int 3
  nopreempt
  unicast_peer {
    192.0.2.12
  }
  virtual_ipaddress {
    192.0.2.1/32 dev ens19
  }
}

vrrp_instance k8s_api_vip {
  state BACKUP
  interface ens19
  virtual_router_id 60
  priority 150
  advert_int 1
  preempt_delay 30
  unicast_peer {
    192.0.2.11
  }
  virtual_ipaddress {
    192.0.2.10/32 dev ens19
  }
}

### openrc-keepalived
#!/sbin/openrc-run
# Generated by routerd. Do not edit by hand.

name='keepalived'
description='routerd managed keepalived VRRP service'
command='/usr/sbin/keepalived'
command_args="'--dont-fork' '--log-console' '--use-file' '/etc/keepalived/keepalived.conf'"
command_background="yes"
pidfile="/run/routerd/openrc/${RC_SVCNAME}.pid"

depend() {
	use net
}

start_pre() {
	checkpath -d -m 0755 /run/routerd/openrc
	checkpath -d -m 0755 '/run/routerd'
	checkpath -d -m 0755 '/var/lib/routerd'
	'/usr/sbin/keepalived' '--config-test' '--use-file' '/etc/keepalived/keepalived.conf'
}

reload() {
	ebegin "Reloading keepalived"
	start-stop-daemon --signal HUP --pidfile "${pidfile}"
	eend $?
}
### openrc-routerd
#!/sbin/openrc-run
# Generated by routerd. Do not edit by hand.

name='routerd'
description='routerd daemon'
command='/usr/local/sbin/routerd'
command_args="'serve'"
command_background="yes"
pidfile="/run/routerd/openrc/${RC_SVCNAME}.pid"

depend() {
	use net
}

start_pre() {
	checkpath -d -m 0755 /run/routerd/openrc
	checkpath -d -m 0755 '/run/routerd'
	checkpath -d -m 0755 '/run/routerd/dhcpv4-client'
	checkpath -d -m 0755 '/run/routerd/dhcpv6-client'
	checkpath -d -m 0755 '/run/routerd/dns-resolver'
	checkpath -d -m 0755 '/run/routerd/pppoe-client'
	checkpath -d -m 0755 '/var/lib/routerd'
	checkpath -d -m 0755 '/var/log/routerd'
	'/usr/local/sbin/routerd' 'check'
}
