# contrib/openrc/routerd
#!/sbin/openrc-run
# shellcheck shell=sh disable=SC2034
# SPDX-License-Identifier: BSD-3-Clause
#
# routerd OpenRC init script for Alpine/OpenRC hosts.
#
# Optional settings may be placed in /etc/conf.d/routerd:
#
# routerd_config="/usr/local/etc/routerd/router.yaml"
# routerd_socket="/run/routerd/routerd.sock"
# routerd_status_socket="/run/routerd/routerd-status.sock"
# routerd_apply_interval="60s"
# routerd_flags=""

name="routerd"
description="declarative router resource applier"

command="${routerd_command:-/usr/local/sbin/routerd}"
pidfile="${routerd_pidfile:-/run/routerd/routerd.pid}"

routerd_config="${routerd_config:-/usr/local/etc/routerd/router.yaml}"
routerd_socket="${routerd_socket:-/run/routerd/routerd.sock}"
routerd_status_socket="${routerd_status_socket:-/run/routerd/routerd-status.sock}"
routerd_apply_interval="${routerd_apply_interval:-60s}"
routerd_flags="${routerd_flags:-}"

command_args="serve --config ${routerd_config} --socket ${routerd_socket} --status-socket ${routerd_status_socket} --apply-interval ${routerd_apply_interval} ${routerd_flags}"

depend() {
    need net
    use logger
}

start_pre() {
    checkpath -d -m 0755 /run/routerd
    checkpath -d -m 0755 /var/lib/routerd
    checkpath -d -m 0755 /var/log/routerd
    if [ ! -x "${command}" ]; then
        eerror "routerd binary ${command} is missing or not executable"
        return 1
    fi
    if [ ! -f "${routerd_config}" ]; then
        eerror "routerd config ${routerd_config} is missing"
        return 1
    fi
}

start() {
    start_pre || return 1
    ebegin "Starting ${name}"
    # shellcheck disable=SC2086
    start-stop-daemon --start --background --make-pidfile --pidfile "${pidfile}" --exec "${command}" -- ${command_args}
    eend $?
}

stop() {
    ebegin "Stopping ${name}"
    start-stop-daemon --stop --pidfile "${pidfile}" --retry TERM/15/KILL/5
    eend $?
}
